# src/ws.ini
[program:websocketserver_publish]
command=python3 /server/publisher.py ; the program (relative uses PATH, can take args)

[program:websocketserver_subscribe]
command=python3 /server/subscriber.py ; the program (relative uses PATH, can take args)
