# src/ws.ini
[program:websocketserver_publish]
command=/bin/bash -c "/root/.bashrc && pyenv global 3.5.2 && python /server/publisher.py"

[program:websocketserver_subscribe]
command=pyenv global 3.5.2 && python /server/subscriber.py
