# src/ws.ini
[program:websocketserver_publish]
command=%(ENV_HOME)s/.pyenv/versions/3.5.2/bin/python /server/publisher.py
user=root

[program:websocketserver_subscribe]
command=%(ENV_HOME)s/.pyenv/versions/3.5.2/bin/python /server/subscriber.py
user=root
